# WebApp/flask.conf
[program:WebApp]
command=gunicorn -c gunicorn.py wsgi:app
directory=/WebApp/
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
