# WebApp/web_nginx.conf
[program:web_nginx]
command=nginx
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=false
stdout_logfile=/WebApp/log/nginx.log
stderr_logfile=/WebApp/log/nginx.err
